# site.conf
[program:site]
command=/home/agoj/venv/bin/uwsgi --ini uwsgi.ini
directory=/home/agoj/site
stopsignal=QUIT
user=agoj
stdout_logfile=/tmp/site.stdout.log
stderr_logfile=/tmp/site.stderr.log
